# roles/zfs/templates/zfs-snapshot.timer
[Unit]
Description=Make ZFS snapshoot every day

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
