# roles/zfs/templates/zfs-snapshot.service
[Unit]
Description=Make ZFS snapshot

[Service]
Type=simple
ExecStart=/bin/bash -c "zfs snapshot -r storage@$(date '+%%Y-%%m-%%d')"
