# roles/restic/templates/restic-backup.timer
[Unit]
Description=Backup via restic every week

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
