# sg-agent/root/crontab
@reboot ssserver -c /root/shadowsocks.json -d start
